# supervisord.conf
[supervisord]
nodaemon	= true

[program:cron]
command		= /usr/sbin/cron -f

[program:httpd]
redirect_stderr=true
command=/usr/local/bin/apache2-foreground
process_name = httpd
